# build/linux/kbs.service
[Unit]
Description=Key Broker Service

[Service]
Type=simple
User=kbs
Group=kbs
ExecStart=/usr/bin/kbs run
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStartSec=0
Restart=on-failure
PermissionsStartOnly=true
RuntimeDirectory=kbs
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
